# deploy/telephony/wwcx-telephony-console.service
[Unit]
Description=WW.CX Big Bird Telephony Operations Console
After=network.target asterisk.service wwcx-numbering-node.service

[Service]
Type=simple
User=wwadmin
Group=wwadmin
WorkingDirectory=/opt/edge1-management-interface
ExecStart=/usr/bin/python3 /opt/edge1-management-interface/server/telephony_status_server.py --host 127.0.0.1 --port 8096
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
