# deploy/systemd/wwcx-mining-operations.timer
[Unit]
Description=Refresh WW.CX mining operations status every minute

[Timer]
OnBootSec=45s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true
Unit=wwcx-mining-operations.service

[Install]
WantedBy=timers.target
