# deploy/systemd/wwcx-operations-version.service
[Unit]
Description=WW.CX Operations Center version exporter

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/python3 /opt/edge1-management-interface/server/operations_version_exporter.py
ReadWritePaths=/var/www/edge1-status
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
