# deploy/systemd/wwcx-operations-timeline.timer
[Unit]
Description=Refresh WW.CX Operations Center timeline

[Timer]
OnBootSec=60s
OnUnitActiveSec=120s
Persistent=true

[Install]
WantedBy=timers.target
