# deploy/bigbird-edge1-connector.service
[Unit]
Description=WW.CX BigBird Edge1 connector lifecycle manager
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=wwadmin
Group=wwadmin
WorkingDirectory=/opt/edge1-management-interface
ExecStart=/usr/bin/python3 /opt/edge1-management-interface/server/bigbird_edge1_connector.py
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/opt/edge1-management-interface
ReadWritePaths=/var/lib/bigbird-edge1-connector
UMask=0077

[Install]
WantedBy=multi-user.target
